# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:php]
command=php5-fpm -R
